# units/sysextmgr-cleanup.service
[Unit]
Description=Daily cleanup of sysext images
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sysextmgrcli cleanup
